# systemd/linuxsup-embedding.service
[Unit]
Description=LinuxSup Face Embedding Service
After=network.target

[Service]
Type=simple
Environment="LD_LIBRARY_PATH=/usr/local/lib/linuxsup"
ExecStart=/usr/local/bin/linuxsup-embedding-service
Restart=always
RestartSec=5

# Security hardening
User=linuxsup
Group=video
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
ReadWritePaths=/run/linuxsup /var/lib/linuxsup
RuntimeDirectory=linuxsup
RuntimeDirectoryMode=0755

# Resource limits
MemoryMax=512M
# CPUQuota=50%  # Disabled - limiting CPU makes face detection too slow

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target